# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

; No `user=` override here: php-fpm's master must start as root (like nginx
; below) so it can reopen the base image's `error_log = /proc/self/fd/2` and
; so its own pool config (www.conf, set to user/group laravel in the
; Dockerfile) can drop worker privileges internally — supervisor pre-dropping
; the master's uid breaks both.
[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.err

[program:nginx]
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log
user=root

; Drains the `jobs` table (QUEUE_CONNECTION=database) — without this, any
; ShouldQueue job (e.g. ProjectImportJob, dispatched when a project's Import
; Setting has run_in_background enabled) is inserted but never processed.
; --tries=3 sends a job to failed_jobs after 3 attempts instead of retrying
; forever; --sleep=3 avoids hammering the DB when the queue is empty.
[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --tries=3 --sleep=3
directory=/app
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/queue-worker.log
stderr_logfile=/var/log/queue-worker.err
user=laravel
