# docker/supervisord.dev.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

; See docker/supervisord.conf for why php-fpm has no `user=` override here.
[program:php-fpm]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.err

[program:nginx]
process_name=%(program_name)s_%(process_num)02d
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log
user=root

; Vite's HMR dev server, so @vite() serves live-reloading assets instead of
; the static public/build the production image uses.
[program:vite]
process_name=%(program_name)s_%(process_num)02d
command=npm run dev -- --host 0.0.0.0
directory=/app
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/vite.log
stderr_logfile=/var/log/vite.err
user=laravel

; Drains the `jobs` table (QUEUE_CONNECTION=database) — without this, any
; ShouldQueue job (e.g. ProjectImportJob, dispatched when a project's Import
; Setting has run_in_background enabled) is inserted but never processed.
; --tries=3 sends a job to failed_jobs after 3 attempts instead of retrying
; forever; --sleep=3 avoids hammering the DB when the queue is empty. Unlike
; php-fpm/vite, a job class edit needs a container restart to be picked up —
; queue:work loads job code once per worker process, not per job.
[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --tries=3 --sleep=3
directory=/app
autostart=true
autorestart=true
numprocs=1
stdout_logfile=/var/log/queue-worker.log
stderr_logfile=/var/log/queue-worker.err
user=laravel
